# templates/restart-xibo.timer.j2
[Unit]
Description=Restart the xibo-linux service daily
Requisite=xibo-linux.service

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
