# templates/restart-xibo.service.j2
[Unit]
Description=Restart xibo-linux service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart xibo-linux.service
